# audiobookshelf/rootfs/etc/services.d/audiobookshelf/run
#!/usr/bin/with-contenv bashio
# Hauptskript (bereits erstellt als run.sh)